# packages/docker/aegis-app-dev/rootfs/etc/s6-overlay/s6-rc.d/aegis-worker/run
#!/bin/sh
cd /app/packages/aegis-server
exec yarn worker:prod
